# systemd/itscape-tracker-update.timer
[Unit]
Description=Run ITScape tracker update twice daily

[Timer]
OnCalendar=*-*-* 06,18:05:00
Persistent=true
RandomizedDelaySec=2m

[Install]
WantedBy=timers.target
