# systemd/itscape-tracker-update.service
[Unit]
Description=Update ITScape tracker via Claude Code

[Service]
Type=oneshot
WorkingDirectory=%h/src/auto-update/itscape
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/src/auto-update/itscape/scripts/auto-update
StandardOutput=journal
StandardError=journal
